# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

; ── Qdrant ───────────────────────────────────────────────────────────

[program:qdrant]
command=/usr/local/bin/qdrant
directory=/workspace/qdrant_data
environment=QDRANT__STORAGE__STORAGE_PATH="/workspace/qdrant_data/storage",QDRANT__SERVICE__HTTP_PORT="6333"
autorestart=true
priority=20
stdout_logfile=/var/log/qdrant.log
stderr_logfile=/var/log/qdrant.err

; ── vLLM — LLM principal (GPU 0) ────────────────────────────────────

[program:vllm_llm]
command=python -m vllm.entrypoints.openai.api_server
    --model /workspace/models/Qwen2.5-Coder-32B-Instruct-AWQ
    --served-model-name Qwen2.5-Coder-32B-Instruct-AWQ
    --quantization awq
    --gpu-memory-utilization 0.92
    --max-model-len 32768
    --tensor-parallel-size 1
    --port 8000
    --api-key %(ENV_API_SECRET_KEY)s
    --enable-auto-tool-choice
    --tool-call-parser hermes
environment=CUDA_VISIBLE_DEVICES="0"
autorestart=true
priority=30
startsecs=30
stdout_logfile=/var/log/vllm_llm.log
stderr_logfile=/var/log/vllm_llm.err

; ── vLLM — Vision (GPU 1) ───────────────────────────────────────────

[program:vllm_vision]
command=python -m vllm.entrypoints.openai.api_server
    --model /workspace/models/Qwen2-VL-7B-Instruct
    --served-model-name Qwen2-VL-7B-Instruct
    --gpu-memory-utilization 0.55
    --max-model-len 8192
    --tensor-parallel-size 1
    --port 8001
    --api-key %(ENV_API_SECRET_KEY)s
environment=CUDA_VISIBLE_DEVICES="1"
autorestart=true
priority=31
startsecs=30
stdout_logfile=/var/log/vllm_vision.log
stderr_logfile=/var/log/vllm_vision.err

; ── Embedding server (GPU 1, shared) ────────────────────────────────

[program:embeddings]
command=python /app/scripts/embedding_server.py
environment=CUDA_VISIBLE_DEVICES="1"
autorestart=true
priority=32
startsecs=45
stdout_logfile=/var/log/embeddings.log
stderr_logfile=/var/log/embeddings.err

; ── Reranker server (GPU 1, shared) ─────────────────────────────────────

[program:reranker]
command=python /app/scripts/reranker_server.py
environment=CUDA_VISIBLE_DEVICES="1"
autorestart=true
priority=33
startsecs=45
stdout_logfile=/var/log/reranker.log
stderr_logfile=/var/log/reranker.err

; ── FastAPI orchestrator ─────────────────────────────────────────────

[program:api]
command=uvicorn src.api.main:app --host 0.0.0.0 --port 8080 --workers 1
directory=/app
autorestart=true
priority=40
startsecs=10
startretries=5
stdout_logfile=/var/log/api.log
stderr_logfile=/var/log/api.err
